# rootfs/etc/services.d/flex-gateway-port-fw/run
#!/usr/bin/env bash
# shellcheck shell=bash

SERVICE=$(basename "$PWD")
CONFIG_FILE="${FLEX_OUTPUT_CONFIG_DIR}/envoy-edge.yaml"


# HEROKU_PRIVATE_IP is defined with Heroku DNS Service Discovery on Private Spaces and it is enabled by default on all applications created after May 31, 2017
if [ -z "${HEROKU_PRIVATE_IP}" ]; then
	export LISTEN_IP=$(hostname -I)
else
	export LISTEN_IP=$(hostname -I | sed s/$HEROKU_PRIVATE_IP// | sed 's/[[:blank:]]//g')
fi

echo "[$SERVICE][info] Starting - envoy port forwarder listening on address: $LISTEN_IP port: $PORT"
exec dockerize -template "envoy-edge.yaml.template:${CONFIG_FILE}" envoy -c "${CONFIG_FILE}" --use-dynamic-base-id -l warn --log-format "[${SERVICE}][%l] %_" --log-path /dev/stdout  --file-flush-interval-msec "${FLEX_SERVICE_ENVOY_FILE_FLUSH_INTERVAL_MS:-1000}"
