# rootfs/etc/services.d/flex-gateway/run
#!/usr/bin/env bash
# shellcheck shell=bash

exec /usr/local/bin/flex-agent
